# systemd-sample.service
[Unit]
Description=Atlas Encryption Service
After=network.target

[Service]
Type=simple
User=lphaap
WorkingDirectory=/home/lphaap/dev/personal/atlas-encryption
ExecStart=/home/lphaap/dev/personal/atlas-encryption/main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
